# configs/system/profile.d/bytefall.sh
export BYTEFALL_RELEASE="0.1 Aurora"
export EDITOR="${EDITOR:-kate}"
export VISUAL="${VISUAL:-kate}"
